# crontab
* * * * * python /task.py >> /ip_checker_logfile
